# systemd/msrc-monitor.service
[Unit]
Description=MSRC Patch Tuesday monitor (monthly)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Replace with the real user at deploy time (USER placeholder).
User=USER
WorkingDirectory=/home/USER/msrc_monitor
# Pushover / healthchecks secrets. Not under git control. Created manually on the Pi at deploy time.
#   e.g.: PUSHOVER_TOKEN=... / PUSHOVER_USER=... / HEALTHCHECK_URL=...
EnvironmentFile=/home/USER/msrc_monitor/env.sh
ExecStart=/home/USER/msrc_monitor/run_monthly.sh

[Install]
WantedBy=multi-user.target
