# systemd/msrc-monitor.timer
[Unit]
Description=MSRC monitor timer (fires the day after Patch Tuesday = the second Wednesday)

[Timer]
# The second Tuesday (Patch Tuesday) falls on the 8th-14th each month. The following day
# (Wednesday) falls on the 9th-15th.
# Any 7-day span from the 9th to the 15th contains exactly one Wednesday, and it is always the
# day after the second Tuesday.
# systemd cannot express "the day after the second Tuesday" directly, so we express it with this
# date range + weekday.
# run_monthly.sh re-checks the same gate (is it the Wednesday within 9..15?) as a second safety check.
OnCalendar=*-*-09..15 Wed 09:00
# If the Pi is asleep and misses the trigger, run once it wakes back up.
Persistent=true

[Install]
WantedBy=timers.target
